# rsscache/bin/crontab
# pwnoogle
0 */2 * * *   root    /htdocs/rsscache/bin/rsscache_cronjob.php
#0 */2 * * *   root    /htdocs/rsscache/bin/rsscache_cronjob.php; /htdocs/rsscache/bin/rsscache_upload.sh
