# rsscache/bin/crontab
# pwnoogle
0 */2 * * *   root    /htdocs/pwnoogle_com/cronjob.sh
