# rsscache/bin/crontab
# pwnoogle
0 */2 * * *   root    /htdocs/rsscache/bin/rsscache_cronjob.sh
#0 */2 * * *   root    /htdocs/rsscache/bin/rsscache_cronjob.sh; /htdocs/rsscache/bin/rsscache_upload.sh
